# build/supervisord.conf
[supervisord]
nodaemon=true

[program:node]
command=node src/index
directory=/usr/src/app
autostart=true
autorestart=true
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:zeronet]
command=bash -c 'sleep 30 && python3 zeronet.py --config_file /etc/zeronet.conf'
directory=/usr/src/app/ZeroNet
autostart=true
autorestart=true
startretries=5
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB

[program:tor]
command=/usr/bin/tor
autostart=true
autorestart=true
startretries=5
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB
